# gui-apps/input-remapper-bin/files/input-remapper.initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

command=/usr/bin/input-remapper-service
command_background=true
pidfile=/run/input-manager.pid

depend() {
	need dbus
}

# stopping the main process leaves hanging forks
stop_post() {
	pkill -f input-remapper-service
}
